# config/php/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.log
user=root

[program:cron]
command=/usr/sbin/cron -f
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron.log
user=root

[include]
files = /etc/supervisor/conf.d/*.conf
